# crontab
1-59 * * * * /bin/bash /home/alix/script
